# webmail/rootfs/services/php/run
#!/bin/sh

echo "Starting PHP..."
exec php-fpm7
